# brains/systemd/robot_runtime.service
[Unit]
Description=Robot Runtime

[Service]
ExecStart=/home/homey/start-stack

[Install]
WantedBy=default.target
WantedBy=multi-user.target
