# brains/systemd/robot_runtime.service
[Unit]
Description=Robot Runtime

[Service]
WorkingDirectory=/home/homey
ExecStart=/usr/bin/docker-compose up

[Install]
WantedBy=default.target
